# context/etc/systemd/system/bootc-fetch-apply-updates.service.d/override.conf
# Don't apply updates just yet
[Service]
ExecStart=
ExecStart=/usr/bin/bootc update --quiet
